# examples/docker/crontab
* * * * * root python /opt/cronio/cronio_worker.py >> /var/log/cron.log 2>&1
